# clamav/clamav-notify.service
[Unit]
Description=ClamAV infection desktop alert

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clamav-notify.sh
